# other/net/openntpd/systemd.d/ntpd.service
[Unit]
Description=ntpd
After=network.target

[Service]
ExecStart=/usr/bin/ntpd -d -s

Nice=-2
IOSchedulingClass=realtime

[Install]
WantedBy=multi-user.target
