# core/filesys/mdadm/systemd.d/mdadm.service
[Unit]
Description=MDADM Event Monitor

[Service]
ExecStart=/usr/bin/mdadm --monitor --scan -y

[Install]
WantedBy=multi-user.target
