# other/graphics/sane-backends/systemd.d/saned@.service
[Unit]
Description=Scanner Service

[Service]
ExecStart=/usr/bin/saned -a usb
StandardInput=socket
